# security_watchdog_15m.cron
*/15 * * * * cd /home/lashaws/.openclaw/workspace/trading_agent_openclaw_v1_1 && ./.venv/bin/python security_watchdog.py >> logs/security_watchdog.log 2>&1
